# systemd/nova-compute.service
[Unit]
Description=Nova Compute server
After=neutron-plugin-openvswitch-agent.service

[Service]
ExecStart=/usr/bin/nova-compute --config-file /etc/nova/nova.conf
User=root
RestartSec=3
Restart=always

[Install]
WantedBy=multi-user.target
